# deploy/cpas.service
[Unit]
Description=CPAs Multi-Site Dashboard
After=network.target nginx.service

[Service]
Type=simple
User=cc-dan
Group=cc-dan
WorkingDirectory=/home/cc-dan/cc/CPAs
EnvironmentFile=/home/cc-dan/cc/CPAs/.env.production
Environment=NODE_ENV=production
ExecStart=/usr/bin/node dist/server/server/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
